# dot_config/systemd/user/backup-snapshot.service
[Unit]
Description=Periodic Backup Snapshot

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=bash -c '"$HOME/.local/share/chezmoi/src/backup/snapshot.sh" -v 2>&1 | logger'

[Install]
WantedBy=default.target
